# macscanner.service
[Unit]
Description=MAC scanner server
Requisite=mysql.service

[Service]
Type=simple
ExecStart=/usr/bin/macscanner.py
User=dan
